# rootfs/opt/nasvcs/runit/lighttpd-log/run
#!/bin/sh
exec tail -f /var/log/lighttpd/access.log
