# rootfs/opt/nasvcs/runit/openssh/run
#!/bin/sh
# TODO: develop a better solution that doesn't pipe to ts
exec /usr/sbin/sshd -D -e 2>&1 | ts '%Y-%m-%dT%H:%M:%S%z sshd:'
